# layout/DEBIAN/postinst
#!/bin/sh
/bin/chown root:wheel /Library/LaunchDaemons/co.cocoanuts.gremlin.plist
/bin/launchctl load /Library/LaunchDaemons/co.cocoanuts.gremlin.plist 2&> /dev/null
exit 0
