# server/layout/Library/LaunchDaemons/co.cocoanuts.gremlin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>2500</integer>
		<key>JetsamPriority</key>
		<integer>-49</integer>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>co.cocoanuts.gremlin</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/gremlind</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>86400</integer>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
